# dpi/files/dpi.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=54
KVER=`uname -r`
KO=dpi-m.ko
start() {
	insmod /lib/modules/${KVER}/${KO}
}

stop() {
	rmmod /lib/modules/${KVER}/${KO}
}
