# luci-mtk/src/themes/dkt/root/etc/uci-defaults/luci-theme-dkt
#!/bin/sh
uci batch <<-EOF
	set luci.themes.DKT=/luci-static/dkt
	set luci.main.mediaurlbase=/luci-static/dkt
        commit luci
EOF
